# osg-ca-certs-updater.cron
# The test is reversed so the expression has an exit status of 0 if the lockfile does not exist.
0 */6 * * * root [ ! -f /var/lock/subsys/osg-ca-certs-updater-cron ] || /usr/sbin/osg-ca-certs-updater -a 23 -x 72 -r 30 -q

